# .yoyo/snapshot/stock-count.service
[Unit]
Description=Stock Count API Service
After=network.target mongodb.service redis.service
Wants=network-online.target

[Service]
Type=notify
User=appuser
Group=appuser
WorkingDirectory=/opt/stock_count
Environment="PATH=/opt/stock_count/.venv/bin"
EnvironmentFile=/opt/stock_count/.env.production

# Main process
ExecStart=/opt/stock_count/.venv/bin/gunicorn backend.server:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --keepalive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --access-logfile /var/log/stock_count/access.log \
    --error-logfile /var/log/stock_count/error.log \
    --log-level info \
    --capture-output \
    --enable-stdio-inheritance

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stock_count /opt/stock_count/backups
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
